# student_files/06/roles/idempotent-demo/files/workshop.service
[Unit]

[Service]
Restart=no
#ExecStart=/usr/bin/python2 /tmp/bin/workshop.py
ExecStart=/usr/bin/python2 /tmp/bin/workshop.py run

[Install]
WantedBy=multi-user.target
